# units/rpi-hw-update.path
[Unit]
Description=Watch for new RPi/Axelera hardware blobs

[Path]
# Watches for new or changed overlays
PathChanged=/lib/firmware/overlays/
Unit=rpi-hw-rebuild.service

[Install]
WantedBy=multi-user.target
